# config/systemd/user/hyprdots-recover-internal-monitor.service
[Unit]
Description=Recover the internal monitor toggle when no external display is connected
Before=graphical-session-pre.target
ConditionPathExists=%h/.local/state/hyprdots/toggles/hypr/internal-monitor-disable.lua

[Service]
Type=oneshot
ExecStart=%h/.local/share/hyprdots/bin/hyprdots-hw-recover-internal-monitor

[Install]
WantedBy=graphical-session-pre.target
